# deploy/systemd/learnkeeper-backup.timer
[Unit]
Description=Run LearnKeeper SQLite backup daily

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
Unit=learnkeeper-backup.service

[Install]
WantedBy=timers.target
